# skel/etc/supervisor/conf.d/300-jenkins.conf
[program:jenkins]
command=%(ENV_SERVICE_JENKINS_CMD)s
user=jenkins
autostart=true
autorestart=unexpected
startsecs=5
startretries=10
priority=300
stopsignal=TERM
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
